# deploy/lightsail/gunicorn.service.example
# Copy to /etc/systemd/system/radhe-cars.service on your Lightsail Ubuntu instance,
# then: sudo systemctl daemon-reload && sudo systemctl enable --now radhe-cars
#
# Adjust User, Group, WorkingDirectory, and EnvironmentFile paths.

[Unit]
Description=Radhe Cars Gunicorn
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/radhe-cars
EnvironmentFile=/var/www/radhe-cars/.env
ExecStart=/var/www/radhe-cars/venv/bin/gunicorn radhe_cars.wsgi:application \
    --bind 127.0.0.1:8000 --workers 2 --timeout 120
Restart=on-failure

[Install]
WantedBy=multi-user.target
